# service/package/statuspanel-service/DEBIAN/postinst
#!/bin/sh

# TODO: Expand the Docker image

systemctl enable statuspanel-service.service
systemctl start statuspanel-service.service
