# service/package/prerm
#!/bin/sh

# Stop and disable the service
systemctl stop statuspanel-service.service
systemctl disable statuspanel-service.service

# Remove the image.
docker image rm --force $IMAGE_SHA
